# deploy/targets/run/supervisord.d/node.ini
[program:nodejs]
command=npm start
directory=/srv/project/code/server
stdout_logfile=/data/log/nodejs.stdout.log
stderr_logfile=/data/log/nodejs.stderr.log
priority = 100
