# deploy/targets/run/supervisord.d/sshd.ini
[program:sshd]
command=/usr/sbin/sshd -D
directory=/usr/sbin
stdout_logfile=/data/log/sshd.stdout.log
stderr_logfile=/data/log/sshd.stderr.log
